# brain/v4/v4_monitor.service
[Unit]
Description=V4 Monitor

[Service]
Type=simple
ExecStart=/usr/bin/v4_monitor

[Install]
WantedBy=multi-user.target
